# com.zeroclaw.hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Always-on launchd agent for the zeroclaw "hermes"/Riff daemon.
  Runs via run.sh so the .env secret injection (Mistral key + Discord token) is applied.
  Install:  cp this to ~/Library/LaunchAgents/ then:
            launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.zeroclaw.hermes.plist
  Stop:     launchctl bootout gui/$(id -u)/com.zeroclaw.hermes
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeroclaw.hermes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sac/Git/llms-local/zeroclaw/run.sh</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sac/Git/llms-local/zeroclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/sac</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sac/Git/llms-local/zeroclaw/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sac/Git/llms-local/zeroclaw/daemon.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
